# Remote BootPicker/PackageRoot/Library/LaunchDaemons/edu.usc.remotebootpicker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>edu.usc.remotebootpicker</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Library/Application Support/USC/remotebp.py</string>
	</array>
	<key>KeepAlive</key>
	<dict>
		<key>PathState</key>
		<dict>
			<key>NetworkState</key>
			<true/>
		</dict>
	</dict>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardOutPath</key>
	<string>/Library/Logs/remotebootpicker.log</string>
	<key>StandardErrorPath</key>
	<string>/Library/Logs/remotebootpicker.log</string>
</dict>
</plist>
