# Remote BootPicker/PackageRoot/Library/LaunchAgents/edu.usc.remotebooter.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>edu.usc.remotebooter</string>
	<key>LimitLoadToSessionType</key>
	<string>LoginWindow</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Library/Application Support/USC/rebootAtLoginWindow.sh</string>
	</array>
	<key>KeepAlive</key>
	<dict>
		<key>PathState</key>
		<dict>
			<key>/Library/Application Support/USC/rebooter/.rebootAtLoginWindow</key>
			<true/>
		</dict>
	</dict>
	<key>RunAtLoad></key>
	<true/>
</dict>
</plist>
